# com.jervis.captainslog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jervis.captainslog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/user/Desktop/BUSSINES/Antigraity/TRADE AI/Jarvis AI/jervis-captains-log.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/user/Desktop/BUSSINES/Antigraity/TRADE AI/Jarvis AI</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/user/Desktop/BUSSINES/Antigraity/TRADE AI/Jarvis AI/state/logs/captainslog-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/user/Desktop/BUSSINES/Antigraity/TRADE AI/Jarvis AI/state/logs/captainslog-launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
